# site-template/sd.service.in
[Unit]
Description=@TARGET_NAME@ Service
Documentation=@DOCURL@
After=network.target 

[Service]
# User=
# Group=

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=@TARGET_NAME@


ExecStart=@JAVAPATH@/java @JAVAOPTS@ -jar @TARGET_FILENAME@ @OPTIONS@ 
#SuccessExitStatus=143 

ExecReload=kill -s TERM $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target
